# ipaddr.service
[Unit]
Description=IP Address Monitor
Before=network-pre.target
Wants=network-pre.target

[Service]
ExecStart=/usr/bin/ipaddr monitor

[Install]
WantedBy=network.target
